# be-python.conf
[program:be-python]
command=/usr/local/bin/gunicorn wsgi:application -c conf-gunicorn.py
directory=/var/www/be-python
autostart=true
autorestart=true